# smartphoton/rootfs/etc/services.d/nodered/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Node-RED with dynamic flows generation using Python
# ==============================================================================

bashio::log.info "Démarrage de Voltronic MQTT..."

# Lire les options de configuration
SERIAL_PORT_1=$(bashio::config 'serial_port_1')
SERIAL_PORT_2=$(bashio::config 'serial_port_2')
SERIAL_PORT_3=$(bashio::config 'serial_port_3')
MQTT_BROKER=$(bashio::config 'mqtt_broker')
MQTT_PORT=$(bashio::config 'mqtt_port')
MQTT_USER=$(bashio::config 'mqtt_user')
MQTT_PASSWORD=$(bashio::config 'mqtt_password')
POLLING_FAST=$(bashio::config 'polling_fast')
POLLING_SLOW=$(bashio::config 'polling_slow')

bashio::log.info "Configuration:"
bashio::log.info "  Serial Port 1: ${SERIAL_PORT_1}"
if [ -n "${SERIAL_PORT_2}" ]; then
    bashio::log.info "  Serial Port 2: ${SERIAL_PORT_2}"
fi
if [ -n "${SERIAL_PORT_3}" ]; then
    bashio::log.info "  Serial Port 3: ${SERIAL_PORT_3}"
fi
bashio::log.info "  MQTT Broker: ${MQTT_BROKER}:${MQTT_PORT}"
bashio::log.info "  Polling rapide: ${POLLING_FAST}s"
bashio::log.info "  Polling lent: ${POLLING_SLOW}s"

# Créer le répertoire s'il n'existe pas
mkdir -p /data/nodered

bashio::log.info "Génération du flows.json..."

# Utiliser Python pour générer le flows.json
python3 << 'PYTHON_SCRIPT'
import json
import os
import sys

# Récupérer les variables d'environnement passées par bash
serial_port_1 = os.environ.get('SERIAL_PORT_1', '')
serial_port_2 = os.environ.get('SERIAL_PORT_2', '')
serial_port_3 = os.environ.get('SERIAL_PORT_3', '')
mqtt_broker = os.environ.get('MQTT_BROKER', 'localhost')
mqtt_port = os.environ.get('MQTT_PORT', '1883')
mqtt_user = os.environ.get('MQTT_USER', '')
mqtt_password = os.environ.get('MQTT_PASSWORD', '')
polling_fast = os.environ.get('POLLING_FAST', '5')
polling_slow = os.environ.get('POLLING_SLOW', '300')

# Lire le template
try:
    with open('/etc/nodered/flows-template.json', 'r') as f:
        flows = json.load(f)
except Exception as e:
    print(f"Erreur lecture template: {e}")
    sys.exit(1)

# Construire le tableau des ports série
serial_ports = [serial_port_1]
if serial_port_2:
    serial_ports.append(serial_port_2)
if serial_port_3:
    serial_ports.append(serial_port_3)

serial_ports_json = json.dumps(serial_ports)

# Fonction récursive pour remplacer les placeholders
def replace_placeholders(obj):
    if isinstance(obj, dict):
        return {k: replace_placeholders(v) for k, v in obj.items()}
    elif isinstance(obj, list):
        return [replace_placeholders(item) for item in obj]
    elif isinstance(obj, str):
        # Remplacer les placeholders
        result = obj
        result = result.replace('{{SERIAL_PORTS}}', serial_ports_json)
        result = result.replace('{{MQTT_BROKER}}', mqtt_broker)
        result = result.replace('{{MQTT_PORT}}', mqtt_port)
        result = result.replace('{{MQTT_USER}}', mqtt_user)
        result = result.replace('{{MQTT_PASSWORD}}', mqtt_password)
        result = result.replace('{{POLLING_FAST}}', polling_fast)
        result = result.replace('{{POLLING_SLOW}}', polling_slow)
        
        # Cas spéciaux pour les serial ports individuels
        if serial_port_1:
            result = result.replace('{{SERIAL_PORT_1}}', serial_port_1)
        if serial_port_2:
            result = result.replace('{{SERIAL_PORT_2}}', serial_port_2)
        if serial_port_3:
            result = result.replace('{{SERIAL_PORT_3}}', serial_port_3)
        
        return result
    else:
        return obj

# Appliquer les remplacements
flows_processed = replace_placeholders(flows)

# Écrire le fichier
try:
    with open('/data/nodered/flows.json', 'w') as f:
        json.dump(flows_processed, f, indent=2)
    print("flows.json généré avec succès")
except Exception as e:
    print(f"Erreur écriture flows.json: {e}")
    sys.exit(1)
PYTHON_SCRIPT

# Vérifier que le fichier a été créé
if [ ! -f /data/nodered/flows.json ]; then
    bashio::log.error "Échec de la génération du flows.json"
    exit 1
fi

bashio::log.info "flows.json généré avec succès"

# Copier le settings.js
cp /etc/nodered/settings.js /data/nodered/settings.js

bashio::log.info "Lancement de Node-RED..."

# Exporter les variables pour Python
export SERIAL_PORT_1
export SERIAL_PORT_2
export SERIAL_PORT_3
export MQTT_BROKER
export MQTT_PORT
export MQTT_USER
export MQTT_PASSWORD
export POLLING_FAST
export POLLING_SLOW

# Démarrer Node-RED
exec node-red \
    --settings /data/nodered/settings.js \
    --userDir /data/nodered
