# smartphoton/rootfs/etc/services.d/nodered/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Node-RED with dynamic flows generation
# ==============================================================================

bashio::log.info "Démarrage de Voltronic MQTT..."

# Lire les options
SERIAL_PORT_1=$(bashio::config 'serial_port_1')
SERIAL_PORT_2=$(bashio::config 'serial_port_2')
SERIAL_PORT_3=$(bashio::config 'serial_port_3')
MQTT_BROKER=$(bashio::config 'mqtt_broker')
MQTT_PORT=$(bashio::config 'mqtt_port')
POLLING_FAST=$(bashio::config 'polling_fast')
POLLING_SLOW=$(bashio::config 'polling_slow')

bashio::log.info "Configuration:"
bashio::log.info "  Serial Port 1: ${SERIAL_PORT_1}"
[ -n "${SERIAL_PORT_2}" ] && bashio::log.info "  Serial Port 2: ${SERIAL_PORT_2}"
[ -n "${SERIAL_PORT_3}" ] && bashio::log.info "  Serial Port 3: ${SERIAL_PORT_3}"
bashio::log.info "  MQTT Broker: ${MQTT_BROKER}:${MQTT_PORT}"

mkdir -p /data/nodered

bashio::log.info "Génération du flows.json..."

# Créer le JSON des ports série
SERIAL_PORTS_JSON="[\"${SERIAL_PORT_1}\""
[ -n "${SERIAL_PORT_2}" ] && SERIAL_PORTS_JSON="${SERIAL_PORTS_JSON},\"${SERIAL_PORT_2}\""
[ -n "${SERIAL_PORT_3}" ] && SERIAL_PORTS_JSON="${SERIAL_PORTS_JSON},\"${SERIAL_PORT_3}\""
SERIAL_PORTS_JSON="${SERIAL_PORTS_JSON}]"

# Remplacer UNIQUEMENT dans les nodes de config (pas dans le code JavaScript)
jq --arg sp1 "${SERIAL_PORT_1}" \
   --arg sp2 "${SERIAL_PORT_2:-/dev/ttyUSB1}" \
   --arg sp3 "${SERIAL_PORT_3:-/dev/ttyUSB2}" \
   --arg broker "${MQTT_BROKER}" \
   --arg port "${MQTT_PORT}" \
   --arg fast "${POLLING_FAST}" \
   --arg slow "${POLLING_SLOW}" \
   --argjson ports "${SERIAL_PORTS_JSON}" \
   '
   # Remplacer dans les serial-port nodes
   map(if .type == "serial-port" then
     if .id == "serial_port_1" then .serialport = $sp1
     elif .id == "serial_port_2" then .serialport = $sp2
     elif .id == "serial_port_3" then .serialport = $sp3
     else . end
   # Remplacer dans mqtt-broker nodes
   elif .type == "mqtt-broker" then
     .broker = $broker | .port = $port
   # Remplacer dans inject nodes (polling)
   elif .type == "inject" and (.name | contains("temps réel")) then
     .repeat = $fast
   elif .type == "inject" and (.name | contains("config")) then
     .repeat = $slow
   else . end)
   ' /etc/nodered/flows-template.json > /data/nodered/flows.json

if [ ! -f /data/nodered/flows.json ]; then
    bashio::log.error "Échec génération flows.json"
    exit 1
fi

cp /etc/nodered/settings.js /data/nodered/settings.js

bashio::log.info "Lancement de Node-RED..."

exec node-red --settings /data/nodered/settings.js --userDir /data/nodered
