# smartphoton/rootfs/etc/services.d/nodered/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Node-RED with config file approach
# ==============================================================================

bashio::log.info "Démarrage de Voltronic MQTT..."

# Lire les options
SERIAL_PORT_1=$(bashio::config 'serial_port_1')
SERIAL_PORT_2=$(bashio::config 'serial_port_2')
SERIAL_PORT_3=$(bashio::config 'serial_port_3')
MQTT_BROKER=$(bashio::config 'mqtt_broker')
MQTT_PORT=$(bashio::config 'mqtt_port')
MQTT_USER=$(bashio::config 'mqtt_user')
MQTT_PASSWORD=$(bashio::config 'mqtt_password')
POLLING_FAST=$(bashio::config 'polling_fast')
POLLING_SLOW=$(bashio::config 'polling_slow')

bashio::log.info "Configuration:"
bashio::log.info "  Serial Port 1: ${SERIAL_PORT_1}"
[ -n "${SERIAL_PORT_2}" ] && bashio::log.info "  Serial Port 2: ${SERIAL_PORT_2}"
[ -n "${SERIAL_PORT_3}" ] && bashio::log.info "  Serial Port 3: ${SERIAL_PORT_3}"
bashio::log.info "  MQTT Broker: ${MQTT_BROKER}:${MQTT_PORT}"

mkdir -p /data/nodered

# Créer le fichier de config pour Node-RED
bashio::log.info "Génération de addon-config.json..."

cat > /data/nodered/addon-config.json << JSONEOF
{
  "serialPorts": [
    "${SERIAL_PORT_1}"$([ -n "${SERIAL_PORT_2}" ] && echo ",
    \"${SERIAL_PORT_2}\"")$([ -n "${SERIAL_PORT_3}" ] && echo ",
    \"${SERIAL_PORT_3}\"")
  ],
  "mqtt": {
    "broker": "${MQTT_BROKER}",
    "port": ${MQTT_PORT},
    "user": "${MQTT_USER}",
    "password": "${MQTT_PASSWORD}"
  },
  "polling": {
    "fast": ${POLLING_FAST},
    "slow": ${POLLING_SLOW}
  }
}
JSONEOF

bashio::log.info "Copie et modification du flows.json..."

# Copier puis modifier avec jq
cp /etc/nodered/flows-template.json /data/nodered/flows.json

jq --arg sp1 "${SERIAL_PORT_1}" \
   --arg sp2 "${SERIAL_PORT_2:-/dev/ttyUSB1}" \
   --arg sp3 "${SERIAL_PORT_3:-/dev/ttyUSB2}" \
   --arg broker "${MQTT_BROKER}" \
   --arg port "${MQTT_PORT}" \
   '
   map(
     if .type == "serial-port" then
       if .id == "serial_port_1" then .serialport = $sp1
       elif .id == "serial_port_2" then .serialport = $sp2
       elif .id == "serial_port_3" then .serialport = $sp3
       else . end
     elif .type == "mqtt-broker" then
       .broker = $broker | .port = $port
     else . end
   )
   ' /data/nodered/flows.json > /data/nodered/flows-temp.json && \
   mv /data/nodered/flows-temp.json /data/nodered/flows.json

cp /etc/nodered/settings.js /data/nodered/settings.js

bashio::log.info "Lancement de Node-RED..."

exec node-red --settings /data/nodered/settings.js --userDir /data/nodered
